# templates/default/carbon-cache.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "carbon-cache v<%= node.graphite.version %>"

start on runlevel [2345]
stop on runlevel [!2345]

env CARBON_CACHE=<%= node.graphite.home %>/bin/carbon-cache.py
env PIDFILE=<%= node.graphite.home %>/storage/carbon-cache.pid

# uncomment this for debugging
# console output

expect daemon
kill timeout 20

pre-stop script
  $CARBON_CACHE --pidfile=$PIDFILE stop

  # don't consider this service stopped until the PID gets removed
  #
  while [ -e $PIDFILE ]; do sleep 1 ; done
  exit 0
end script

exec $CARBON_CACHE --pidfile=$PIDFILE start
