# templates/default/statsite.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "statsite"

start on runlevel [2345]
stop on runlevel [!2345]

# uncomment this for debugging
# console output

exec statsite -c <%= node.graphite.home %>/conf/statsite.conf
