# templates/default/metrics/metric.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "<%= @name %> metric collector"

start on runlevel [2345]
stop on runlevel [!2345]

console output

script
  <%= @command %>
end script
